# roles/channelpedia/templates/channelpedia.timer.j2
[Unit]
Description=Upload channels.conf to channelpedia.yavdr.com

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
