# roles/yavdr-xorg/templates/user@666.service.d/x-dependency.conf.j2
[Unit]
Wants=x@vt7.service
After=x@vt7.service

[Service]
TimeoutStopSec=20
KillMode=mixed
