# roles/yavdr-desktop/templates/systemd/user/detect-second-display.service.j2
[Unit]
Description=Detect second DISPLAY using xrandr

[Service]
Type=oneshot
ExecStart=%h/bin/detect-second-display

[Install]
WantedBy=yavdr-desktop.target
