# roles/pulseaudio/templates/dbus_pulsectl/dbus-pulsectl.service.j2
[Unit]
Description=Control pulseaudio output via DBus

[Service]
Type=simple
ExecStart=/var/lib/vdr/bin/dbus_pulsectl
KillSignal=SIGINT

[Install]
WantedBy=yavdr-desktop.target
