# roles/yavdr-xorg/templates/systemd/user/irexec.service.j2
[Unit]
Description=LIRC command handler

[Service]
Type=simple
ExecStart=/usr/bin/irexec %h/.lircrc

[Install]
WantedBy=yavdr-desktop.target
