# roles/yavdr-xorg/templates/systemd/user/tmux.service.j2
[Unit]
Description=Start tmux in detached session

[Service]
EnvironmentFile=-%h/.session-env
Type=forking
ExecStart=/usr/bin/tmux new-session -s use-session -d
ExecStop=/usr/bin/tmux kill-session -t %u
Restart=always

[Install]
WantedBy=default.target
