# roles/yavdr-xorg/templates/systemd/user/wm-exit.service.j2
[Unit]
Description=exit window manager gracefully

[Service]
ExecStart=/bin/true
ExecStop=/usr/bin/openbox --exit
RemainAfterExit=True

[Install]
WantedBy=default.target
