# roles/install-sundtek/templates/systemd/vdr.service.d/sundtek.conf.j2
{{ ansible_managed | comment }}

{% if install_sundtek %}
[Service]
Environment=LD_PRELOAD=/opt/lib/libmediaclient.so
{% endif %}
