# roles/install-sundtek/templates/udev/80-mediasrv-sundtek.rules.j2
{{ ansible_managed | comment }}

{% if install_sundtek %}
SUBSYSTEM=="usb", ACTION=="add", RUN+="/opt/bin/mediaclient --systemdcheck"
{% endif %}
