# roles/udiskie/templates/udiskie.service.j2
[Unit]
Description=dynamic mounts for the user session

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
