# roles/yavdr-remote/templates/lircd.service.d/lircd2uinput.conf.j2
{{ ansible_managed | comment }}
[Service]
ExecStartPost=/usr/bin/lircd2uinput-add /var/run/lirc/lircd0
ExecStopPost=/usr/bin/lircd2uinput-remove /var/run/lirc/lircd0
