# roles/yavdr-remote/templates/lircd.service.d/lircd2uinput.conf.j2
[Service]
ExecStartPre=/usr/bin/lircd2uinput-add /var/run/lirc/lircd0
ExecStartPre=/usr/bin/lircd2uinput-remove /var/run/lirc/lircd0
